# lisa/lisa_client_initd

#! /bin/sh
### BEGIN INIT INFO
# Provides:          lisa_client
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start/stop Lisa client
### END INIT INFO

logger "Lisa Server: Start script executed"
LISA_CLIENT_PATH="/home/julien/Sandbox/LISA-CLIENT-Linux/"
export PYTHONPATH="$LISA_CLIENT_PATH:$PYTHONPATH"

case "$1" in
  start)
    logger "Lisa Client: Starting"
    echo "Starting Lisa Client..."
 cd "$LISA_CLIENT_PATH/"
 twistd -y lisa.tac --logfile="$LISA_CLIENT_PATH/lisa.log" --pidfile="$LISA_CLIENT_PATH/Lisa_client.pid"
    ;;
  stop)
    logger "Lisa Client: Stopping"
    echo "Stopping Lisa client..."
    kill `cat $LISA_CLIENT_PATH/Lisa_client.pid`
    ;;
  *)
    logger "Lisa client: Invalid usage"
    echo "Usage: /etc/init.d/lisa_client {start|stop}"
    exit 1
    ;;
esac

